# infra/cron/evals.crontab
# Phase 4c eval schedule — host crontab template.
#
# Run the eval CLI every night at 03:17 UTC against the configured tenant.
# Non-zero exit (RAG recall below configured gates or injection-resistance
# below 100% by default) causes cron to mail the script's stderr to MAILTO —
# wire that up to your ops inbox, PagerDuty email integration, or whichever
# channel you actually read. Don't rely on a Slack hook here; if the eval
# pipeline is down the hook may be too.
#
# Why 03:17 instead of 03:00: avoids the synchronized-cron-at-the-top-of-the-
# hour spike that hits OpenAI's rate limiter from every project running its
# nightly job at the same time.
#
# Install on the host:
#   cd /opt/fielddesk/apps/worker
#   python3 -m venv .venv
#   . .venv/bin/activate
#   pip install -e .
#   crontab -e
#   # paste the line below, replacing FIELDDESK_HOME and TENANT
#
# Or as a system-wide /etc/cron.d/fielddesk-evals:
#   # /etc/cron.d/fielddesk-evals
#   17 3 * * * fielddesk /opt/fielddesk/scripts/eval.sh <tenant-uuid>

MAILTO=ops@example.com
FIELDDESK_HOME=/opt/fielddesk

# m  h dom mon dow  command
17   3  *   *   *   cd "$FIELDDESK_HOME" && ./scripts/eval.sh "<tenant-uuid>" all
